# containers/rsync/crontab
50 * * * * rsync -rtlvH --delete-after --delay-updates --safe-links rsync://de3.mirror.archlinuxarm.org/archlinux-arm /share/archlinuxarm
